# digital-misfit-network-services/fs/etc/cron.d/cron-acme-issue
*/5 * * * * root . $HOME/.cron.env; /root/ssl/acme-issue.sh >/proc/1/fd/1 2>/proc/1/fd/2
# Blank line for valid crontab
